# tegrity/tegrity-firstboot.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

# Tegrity First Boot Service template
[Unit]
After={after}
Before={before}
ConditionFirstBoot=yes
ConditionPathIsReadWrite=/etc
Conflicts=shutdown.target
Description="Tegrity First Boot Service"

[Service]
ExecStart={execstart}
RemainAfterExit=yes
Type=oneshot

[Install]
WantedBy=sysinit.target